# stosos.service
[Unit]
Description=StosOS Desktop Environment
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=STOSOS_USER_PLACEHOLDER
Group=STOSOS_USER_PLACEHOLDER
WorkingDirectory=STOSOS_DIR_PLACEHOLDER
Environment=DISPLAY=:0
Environment=PULSE_RUNTIME_PATH=/run/user/UID_PLACEHOLDER/pulse
ExecStart=STOSOS_DIR_PLACEHOLDER/venv/bin/python STOSOS_DIR_PLACEHOLDER/main.py
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=1G
CPUQuota=80%

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=STOSOS_DIR_PLACEHOLDER/data STOSOS_DIR_PLACEHOLDER/logs STOSOS_DIR_PLACEHOLDER/config

[Install]
WantedBy=graphical-session.target